# roles/openresty/templates/openresty.service.j2
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/rc.d/init.d/openresty
Description=SYSV: OpenResty is a scalable web platform by extending
Before=runlevel3.target
Before=runlevel4.target
Before=runlevel5.target
Before=shutdown.target
After=network-online.target
After=redis_6380.service
After=network.service
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=no
PIDFile=/usr/local/openresty/nginx/logs/nginx.pid
ExecStart=/etc/rc.d/init.d/openresty start
ExecStop=/etc/rc.d/init.d/openresty stop
ExecReload=/etc/rc.d/init.d/openresty reload
